# image/pi-gen-config/stage-eigsep/00-eigsep-install/files/systemd/eigsep-first-boot.service
[Unit]
Description=EIGSEP first-boot role selector
# Applies the operator's eigsep-role.conf once per image flash: enables
# the matching role's [services.*] entries, then self-disables. Idempotent —
# editing the role file and `systemctl enable eigsep-first-boot.service;
# reboot` re-rolls the Pi.
#
# The role conf lives on the SD card's FAT boot partition, which Trixie's
# pi-gen mounts at /boot/firmware.
ConditionPathExists=/boot/firmware/eigsep-role.conf
After=multi-user.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/eigsep-field _apply-role
RemainAfterExit=no

[Install]
WantedBy=multi-user.target
