# image/pi-gen-config/stage-eigsep/00-eigsep-install/files/systemd/eigsep-panda.service
[Unit]
Description=EIGSEP PandaClient (antenna-side producer)
After=network-online.target redis-server.service
Wants=redis-server.service
PartOf=eigsep-panda.target

[Service]
Type=simple
User=eigsep
ExecStart=/opt/eigsep/venv/bin/python -m eigsep_observing.client
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=eigsep-panda.target
