# image/pi-gen-config/stage-eigsep/00-eigsep-install/files/systemd/eigsep-panda.target
[Unit]
Description=EIGSEP panda role plumbing (antenna-side Pico + VNA bridges; operator launches panda_observe manually)
Wants=redis-server.service picomanager.service cmtvna.service
After=network-online.target

[Install]
WantedBy=multi-user.target
