# image/pi-gen-config/stage-eigsep/files/systemd/eigsep-dhcp.target
[Unit]
Description=EIGSEP DHCP server (exactly one Pi per LAN)
Wants=isc-dhcp-server.service
After=network-online.target

[Install]
WantedBy=multi-user.target
