# image/pi-gen-config/stage-eigsep/files/systemd/eigsep-backend.target
[Unit]
Description=EIGSEP backend role (digital backend: SNAP + observer)
Wants=redis-server.service eigsep-observe.service eigsep-observe-writer.service
After=network-online.target

[Install]
WantedBy=multi-user.target
