# image/pi-gen-config/stage-eigsep/files/systemd/eigsep-panda.target
[Unit]
Description=EIGSEP panda role (antenna-side: Pico + VNA + PandaClient)
Wants=redis-server.service eigsep-panda.service picomanager.service cmtvna.service
After=network-online.target

[Install]
WantedBy=multi-user.target
