# image/pi-gen-config/stage-eigsep/00-eigsep-install/files/systemd/chrony-wait.service.d/eigsep.conf
# Debian's chrony-wait.service ships without an [Install] section, so
# `systemctl enable` is a no-op: nothing pulls the gate into the boot
# graph and any `After=chrony-wait.service` silently doesn't wait. This
# drop-in adds the WantedBy so the `enable-always` step in
# _image_install.py actually wires the gate into multi-user.target on
# every Pi at image build time.
[Install]
WantedBy=multi-user.target
